# ops/production/chathtml-backup.timer
[Unit]
Description=Nightly ChatHTML backup

[Timer]
OnCalendar=*-*-* 02:20:00
RandomizedDelaySec=20m
Persistent=true

[Install]
WantedBy=timers.target
